# deploy/smartsort.service
[Unit]
Description=SmartSort - Organizador de Ficheiros Inteligente com IA
After=network.target

[Service]
# O instalador vai substituir estes placeholders pelo seu utilizador e caminho reais
User=USER_PLACEHOLDER
WorkingDirectory=WORKING_DIR_PLACEHOLDER
ExecStart=WORKING_DIR_PLACEHOLDER/venv/bin/python -u -m smartsort
Environment=PYTHONPATH=WORKING_DIR_PLACEHOLDER/src
Environment=PYTHONUNBUFFERED=1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
